# naming/core/crontab
* * * * * bash -c 'python $WORKDIR/manage.py check_storage_servers_health' &> /tmp/cron.log
