# naming/core/crontab
* * * * * python $(WORKDIR)/manage.py check_storage_servers_health > /tmp/naming.log
